# megacmd/rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting terminal..."
cd /root
exec ttyd -p 7681 -W -d 0 /bin/bash -l
