# megacmd/rootfs/etc/s6-overlay/s6-rc.d/mqtt-listener/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

if ! bashio::services.available "mqtt"; then
  exit 0
fi

MQTT_HOST=$(bashio::services mqtt "host")
MQTT_PORT=$(bashio::services mqtt "port")
MQTT_USER=$(bashio::services mqtt "username")
MQTT_PASS=$(bashio::services mqtt "password")

mosquitto_sub -h "$MQTT_HOST" -p "$MQTT_PORT" -u "$MQTT_USER" -P "$MQTT_PASS" -t "megacmd/command" | while read -r cmd; do
  if ! mega-whoami >/dev/null 2>&1; then
    continue
  fi

  case "$cmd" in
  clean_local_debris)
    bashio::log.info "Clean local debris..."

    # Use for loop, because direct call of find caused errors
    FOLDERS=(/media /homeassistant /config /backup /share)
    for folder in "${FOLDERS[@]}"; do
      if [ -d "$folder" ]; then
        find "$folder" -type d -name ".debris" -exec rm -rf {} + 2>/dev/null
      fi
    done

    /usr/local/bin/publish_mega_size.sh
    ;;

  clean_cloud_debris)
    bashio::log.info "Clean cloud debris..."
    mega-rm -r -f //bin/SyncDebris
    mega-mkdir //bin/SyncDebris
    /usr/local/bin/publish_mega_size.sh
    ;;

  *)
    bashio::log.info "Unknown command: $cmd"
    ;;

  esac

done
